# cron.d/ntpdate
 *    */4  *    *    *          root    /usr/sbin/ntpdate-debian
